# scripts/urine-pi.service
# Reference template — actual service is installed by scripts/setup-pi.sh
# Placeholders are filled in at install time with real paths.
#
# Manual install:
#   sudo cp scripts/urine-pi.service /etc/systemd/system/
#   # Edit User=, WorkingDirectory=, EnvironmentFile=, ExecStart= manually
#   sudo systemctl daemon-reload && sudo systemctl enable --now urine-pi

[Unit]
Description=Urine Analyzer — Pi Client
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=120
StartLimitBurst=5

[Service]
Type=simple
User=PLACEHOLDER_USER
WorkingDirectory=PLACEHOLDER_DIR
EnvironmentFile=PLACEHOLDER_DIR/.env
ExecStart=PLACEHOLDER_UV run python main.py --role client
Restart=always
RestartSec=15
KillMode=process
StandardOutput=journal
StandardError=journal
SyslogIdentifier=urine-pi

[Install]
WantedBy=multi-user.target
